# deploy/services.d/gunicorn/run
#!/bin/execlineb -P
gunicorn -c /app/deploy/gunicorn.py www.wsgi:application
